# root/etc/s6-overlay/s6-rc.d/rsyslog/run
#!/usr/bin/with-contenv bash

sleep 15
exec rsyslogd -n