# root/etc/s6-overlay/s6-rc.d/cron/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "[cron] Starting cron daemon to handle certificate renewals..." | ts '%Y-%m-%d %H:%M:%S'

# Create required directories
mkdir -p /etc/cron.d /var/spool/cron/crontabs /etc/cron.daily /etc/cron.hourly /etc/cron.monthly /etc/cron.weekly

# Ensure cron log exists
touch /var/log/cron.log
chmod 640 /var/log/cron.log

# Fix permissions for cron directories
chmod -R 0644 /etc/cron.d
chmod -R 0755 /etc/cron.daily /etc/cron.hourly /etc/cron.monthly /etc/cron.weekly

# Create /etc/passwd entry for acme user if needed
if ! grep -q "^acme:" /etc/passwd; then
    echo "acme:x:1000:1000:acme:/config/acme:/bin/sh" >> /etc/passwd
fi

# Create /etc/group entry for acme group if needed
if ! grep -q "^acme:" /etc/group; then
    echo "acme:x:1000:" >> /etc/group
fi

# Check for existing renewal job
if [ -f /etc/cron.d/acme-renewal ]; then
    echo "[cron] Found existing acme-renewal job" | ts '%Y-%m-%d %H:%M:%S'
else
    # Create a default renewal job if none exists
    echo "[cron] Creating default renewal job" | ts '%Y-%m-%d %H:%M:%S'
    cat << EOF > /etc/cron.d/acme-renewal
# Run certificate renewal at 2:30 AM on Monday and Thursday
30 2 * * 1,4 acme /usr/local/bin/renew-certs.sh > /var/log/acme-renewals.log 2>&1
EOF
    chmod 0644 /etc/cron.d/acme-renewal
fi

# Create a helper script to manually run certificate renewal
cat << 'EOF' > /usr/local/bin/run-acme-renewal
#!/bin/bash
echo "Running manual certificate renewal..."
if [ -f /usr/local/bin/renew-certs.sh ]; then
    s6-setuidgid acme /usr/local/bin/renew-certs.sh
    echo "Certificate renewal complete."
else
    echo "Error: Renewal script not found at /usr/local/bin/renew-certs.sh"
    exit 1
fi
EOF
chmod +x /usr/local/bin/run-acme-renewal

echo "[cron] Starting cron daemon in foreground mode" | ts '%Y-%m-%d %H:%M:%S'
# Start cron daemon with debugging
exec crond -f -L /var/log/cron.log 