# root/etc/s6-overlay/s6-rc.d/cron/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "[cron] starting.." | ts '%Y-%m-%d %H:%M:%S'

exec crond -f -d 0 >/dev/null 2>&1
