# root/etc/s6-overlay/s6-rc.d/haproxy/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

file1="/config/acme/ca/thumbprint"
file2="/config/acme/tls1-params/ffdhe2048"
HAPROXY_YAML="/config/haproxy.yaml"

source /scripts/debug.sh;

# Loop until both files exist
while [ ! -f "$file1" ] || [ ! -f "$file2" ]; do
    echo "[Haproxy] waiting for $file1 and $file2 to be created before starting..." | ts '%Y-%m-%d %H:%M:%S'
    sleep 3
done

if [ "${HA_DEBUG_ENABLED}" == "true" ]; then
    s6-setuidgid haproxy haproxy -vv
else
    s6-setuidgid haproxy haproxy -v
fi

# Ensure the socket files exist and have the correct permissions
for sock in frontend-offloading.sock frontend-offloading-ip-protection.sock; do
    if [ ! -f "/var/lib/haproxy/$sock" ]; then
        touch "/var/lib/haproxy/$sock"
        chmod 660 "/var/lib/haproxy/$sock"
        chown haproxy:haproxy "/var/lib/haproxy/$sock"
    fi
done

# Start HAProxy if it's not already running
if ! pgrep -x "haproxy" > /dev/null; then
    echo "[Haproxy] started!" | ts '%Y-%m-%d %H:%M:%S'

    if [ "${HA_DEBUG_ENABLED}" == "true" ]; then
        s6-setuidgid haproxy haproxy -f -d /config/haproxy.cfg
    else
        s6-setuidgid haproxy haproxy -f /config/haproxy.cfg
    fi
else
    echo "[Haproxy] is already running." | ts '%Y-%m-%d %H:%M:%S'
    exit 0
fi