# root/etc/s6-overlay/s6-rc.d/port-forwarding/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# This script sets up port forwarding for HTTP/3 when MIXED_SSL_MODE is enabled
# It forwards UDP traffic from port 443 to port 8443 inside the container using socat
# This runs as a oneshot service during container startup

if [ "${MIXED_SSL_MODE}" = "true" ]; then
    echo "[port-forwarding] MIXED_SSL_MODE is enabled, setting up UDP port forwarding from 443 to 8443..." | ts '%Y-%m-%d %H:%M:%S'

    # Check if socat is available
    if ! command -v socat &> /dev/null; then
        echo "[port-forwarding] Error: socat is not installed. Cannot set up port forwarding." | ts '%Y-%m-%d %H:%M:%S'
        exit 0
    fi

    # Kill any existing socat processes for this port
    pkill -f "socat UDP-LISTEN:443" 2>/dev/null

    # Start socat in the background to forward UDP traffic
    socat UDP-LISTEN:443,fork UDP:127.0.0.1:8443 &
    SOCAT_PID=$!

    # Check if socat started successfully
    if [ -n "$SOCAT_PID" ] && kill -0 $SOCAT_PID 2>/dev/null; then
        echo "[port-forwarding] Successfully started socat UDP port forwarding from 443 to 8443 (PID: $SOCAT_PID)" | ts '%Y-%m-%d %H:%M:%S'
        
        # Save the PID to a file so we can stop it later if needed
        echo $SOCAT_PID > /var/run/socat-port-forward.pid
    else
        echo "[port-forwarding] Failed to start socat. HTTP/3 may not work correctly." | ts '%Y-%m-%d %H:%M:%S'
    fi
else
    echo "[port-forwarding] MIXED_SSL_MODE is disabled, no port forwarding needed" | ts '%Y-%m-%d %H:%M:%S'
fi