# root/etc/s6-overlay/s6-rc.d/port-forwarding/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# This script sets up port forwarding for HTTP/3 when MIXED_SSL_MODE is enabled
# It forwards UDP traffic from port 443 to port 8443 inside the container
# This must run as root, which is why it's a separate S6 service

if [ "${MIXED_SSL_MODE}" = "true" ]; then
    echo "[port-forwarding] MIXED_SSL_MODE is enabled, setting up UDP port forwarding from 443 to 8443..." | ts '%Y-%m-%d %H:%M:%S'

    # Check if iptables is available
    if ! command -v iptables &> /dev/null; then
        echo "[port-forwarding] Warning: iptables is not installed. HTTP/3 may not work correctly with MIXED_SSL_MODE." | ts '%Y-%m-%d %H:%M:%S'
        echo "[port-forwarding] Continuing without port forwarding..." | ts '%Y-%m-%d %H:%M:%S'
        exit 0
    fi

    # Check if the rule already exists to avoid duplicates
    if ! iptables -t nat -C PREROUTING -p udp --dport 443 -j REDIRECT --to-port 8443 2>/dev/null; then
        # Add the port forwarding rule
        if iptables -t nat -A PREROUTING -p udp --dport 443 -j REDIRECT --to-port 8443; then
            echo "[port-forwarding] Successfully set up UDP port forwarding from 443 to 8443" | ts '%Y-%m-%d %H:%M:%S'
        else
            echo "[port-forwarding] Failed to set up port forwarding. HTTP/3 may not work correctly." | ts '%Y-%m-%d %H:%M:%S'
        fi
    else
        echo "[port-forwarding] Port forwarding rule already exists" | ts '%Y-%m-%d %H:%M:%S'
    fi

    # List current rules for verification
    echo "[port-forwarding] Current NAT rules:" | ts '%Y-%m-%d %H:%M:%S'
    iptables -t nat -L -n -v | grep -i redirect | ts '%Y-%m-%d %H:%M:%S'
else
    echo "[port-forwarding] MIXED_SSL_MODE is disabled, no port forwarding needed" | ts '%Y-%m-%d %H:%M:%S'
fi