# root/etc/s6-overlay/s6-rc.d/port-forwarding/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# This script sets up port forwarding for HTTP/3 when MIXED_SSL_MODE is enabled
# It forwards UDP traffic from port 443 to port 8443 inside the container using socat
# This runs as a longrun service during container lifetime

echo "[port-forwarding] Checking environment..." | ts '%Y-%m-%d %H:%M:%S'
echo "[port-forwarding] MIXED_SSL_MODE=${MIXED_SSL_MODE}" | ts '%Y-%m-%d %H:%M:%S'

# Check if we're running as root
if [ "$(id -u)" != "0" ]; then
    echo "[port-forwarding] Warning: Not running as root, this may cause permission issues" | ts '%Y-%m-%d %H:%M:%S'
fi

# Check for existing listeners on port 443
echo "[port-forwarding] Current UDP listeners on port 443:" | ts '%Y-%m-%d %H:%M:%S'
netstat -tulpn | grep 443 | grep udp | ts '%Y-%m-%d %H:%M:%S' || echo "[port-forwarding] No UDP listeners found on port 443" | ts '%Y-%m-%d %H:%M:%S'

if [ "${MIXED_SSL_MODE}" = "true" ]; then
    echo "[port-forwarding] MIXED_SSL_MODE is enabled, setting up UDP port forwarding from 443 to 8443..." | ts '%Y-%m-%d %H:%M:%S'

    # Check if socat is available
    if ! command -v socat &> /dev/null; then
        echo "[port-forwarding] Error: socat is not installed. Cannot set up port forwarding." | ts '%Y-%m-%d %H:%M:%S'
        exit 1
    fi

    # Kill any existing socat processes for this port
    echo "[port-forwarding] Killing any existing socat processes..." | ts '%Y-%m-%d %H:%M:%S'
    pkill -f "socat UDP-LISTEN:443" 2>/dev/null

    # Check for iptables rules and remove them if they exist
    if command -v iptables &> /dev/null; then
        echo "[port-forwarding] Checking for existing iptables rules..." | ts '%Y-%m-%d %H:%M:%S'
        if iptables -t nat -C PREROUTING -p udp --dport 443 -j REDIRECT --to-port 8443 2>/dev/null; then
            echo "[port-forwarding] Removing existing iptables rule..." | ts '%Y-%m-%d %H:%M:%S'
            iptables -t nat -D PREROUTING -p udp --dport 443 -j REDIRECT --to-port 8443
        fi
    fi

    # Start socat with more verbose options
    echo "[port-forwarding] Starting socat UDP port forwarding from 443 to 8443" | ts '%Y-%m-%d %H:%M:%S'
    exec socat -d -d UDP-LISTEN:443,fork,reuseaddr UDP:127.0.0.1:8443
else
    echo "[port-forwarding] MIXED_SSL_MODE is disabled, no port forwarding needed" | ts '%Y-%m-%d %H:%M:%S'
    # Sleep forever since this is a longrun service
    exec sleep infinity
fi