# root/etc/s6-overlay/s6-rc.d/haproxy-config/finish
#!/usr/bin/with-contenv bash

echo "[Haproxy] configuration script completed" | ts '%Y-%m-%d %H:%M:%S'