# root/etc/s6-overlay/s6-rc.d/acme/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

source /scripts/acme_functions.sh # issue_cert and deploy_cert functions
source /scripts/debug.sh

ACME_THUMBPRINT="/config/acme/ca/thumbprint"
HAPROXY_PID_FILE="/var/run/haproxy/haproxy.pid"
HAPROXY_SOCKET="/var/lib/haproxy/admin.sock"

# Loop until acme thumbprint file exists
while [ ! -f "$ACME_THUMBPRINT" ]; do
    echo "[acme] waiting for $ACME_THUMBPRINT to be created before setting up acme..." | ts '%Y-%m-%d %H:%M:%S'
    sleep 3
done

# Wait for HAProxy to be fully started by checking for its PID file and socket
echo "[acme] waiting for HAProxy to be fully started..." | ts '%Y-%m-%d %H:%M:%S'
while [ ! -f "$HAPROXY_PID_FILE" ] || [ ! -S "$HAPROXY_SOCKET" ]; do
    echo "[acme] waiting for HAProxy to be fully started (PID file: $([ -f "$HAPROXY_PID_FILE" ] && echo "exists" || echo "missing"), Socket: $([ -S "$HAPROXY_SOCKET" ] && echo "exists" || echo "missing"))..." | ts '%Y-%m-%d %H:%M:%S'
    sleep 5
done

# Additional check to ensure HAProxy is responding
MAX_ATTEMPTS=12
ATTEMPT=0
while [ $ATTEMPT -lt $MAX_ATTEMPTS ]; do
    if echo "show info" | socat stdio "unix-connect:$HAPROXY_SOCKET" &>/dev/null; then
        echo "[acme] HAProxy is fully operational, proceeding with certificate checks..." | ts '%Y-%m-%d %H:%M:%S'
        break
    fi
    echo "[acme] HAProxy is not responding yet, waiting (attempt $((ATTEMPT+1))/$MAX_ATTEMPTS)..." | ts '%Y-%m-%d %H:%M:%S'
    ATTEMPT=$((ATTEMPT+1))
    sleep 5
done

# Process certificates
check_for_missing_domain_certs "no"

# Make sure renewal script exists and is executable
echo "[acme] Setting up certificate renewal job..." | ts '%Y-%m-%d %H:%M:%S'

if [ ! -x "/usr/local/bin/renew-certs.sh" ]; then
    setup_acme_renewal
else
    # Ensure script has correct permissions
    chmod +x /usr/local/bin/renew-certs.sh
    chown "${USER}:${USER}" /usr/local/bin/renew-certs.sh
fi

# Verify that cron daemon is running (should be started by the cron service)
if pgrep -x "crond" > /dev/null; then
    echo "[acme] Cron daemon is running" | ts '%Y-%m-%d %H:%M:%S'
else
    echo "[acme] Warning: Cron daemon not detected, certificates may not auto-renew" | ts '%Y-%m-%d %H:%M:%S'
fi

# Check for cron.d file
if [ -f "/etc/cron.d/acme-renewal" ]; then
    echo "[acme] Found cron job configuration in /etc/cron.d/acme-renewal" | ts '%Y-%m-%d %H:%M:%S'
    cat "/etc/cron.d/acme-renewal" | ts '%Y-%m-%d %H:%M:%S'
else 
    echo "[acme] Warning: No cron job found in /etc/cron.d/acme-renewal" | ts '%Y-%m-%d %H:%M:%S'
    setup_acme_renewal
fi

echo "[acme] Certificate renewal setup complete" | ts '%Y-%m-%d %H:%M:%S'

# Keep the service running
tail -f /dev/null